# etc/cron.hourly/droplet-agent
#!/bin/sh
/bin/bash /opt/digitalocean/droplet-agent/scripts/update.sh >/var/log/droplet-agent.update.log 2>&1
